# etc/supervisor/conf.d/api.conf
[program:api]
command=/var/site/projects/nodejs/api/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api-error.log
