# etc/supervisor/conf.d/chat.conf
[program:chat]
command=/var/site/projects/nodejs/chat/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/chat.log
stderr_logfile=/var/log/supervisor/chat-error.log
