# etc/supervisor/conf.d/terminal.conf
[program:wetty]
command=/var/site/projects/nodejs/wetty/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/wetty.log
stderr_logfile=/var/log/supervisor/wetty-error.log
