# etc/supervisor/conf.d/twig.conf
[program:twig]
command=/var/site/projects/node/twig-express/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/twig.log
stderr_logfile=/var/log/supervisor/twig-error.log
