# etc/supervisor/conf.d/wetty.conf
[program:wetty]
command=/home/site/projects/node/wetty/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/wetty.log
stderr_logfile=/var/log/supervisor/wetty-error.log
